# supervisor-chirpstack.conf
; Servizi gestiti da supervisord (avvio e restart automatico, niente screen)

; MQTT per ChirpStack (deve essere il primo)
[program:mosquitto]
command=mosquitto -c /etc/mosquitto/mosquitto.conf
priority=100
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/mosquitto.log
stderr_logfile=/var/log/supervisor/mosquitto_err.log

; ChirpStack Network Server
[program:chirpstack-network-server]
command=chirpstack-network-server
priority=200
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/network-server.log
stderr_logfile=/var/log/supervisor/network-server_err.log

; ChirpStack Application Server (UI + API)
[program:chirpstack-application-server]
command=chirpstack-application-server
priority=200
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/application-server.log
stderr_logfile=/var/log/supervisor/application-server_err.log

; ChirpStack Gateway Bridge (UDP 1700 -> MQTT), dopo MQTT
[program:chirpstack-gateway-bridge]
command=chirpstack-gateway-bridge -c /etc/chirpstack-gateway-bridge/chirpstack-gateway-bridge.toml
priority=300
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/var/log/supervisor/gateway-bridge.log
stderr_logfile=/var/log/supervisor/gateway-bridge_err.log

; LWN-Simulator (porta 9000): avvio ritardato via wrapper, restart automatico da supervisord
[program:lwn-simulator]
command=/root/run_lwn_delayed.sh
priority=400
autostart=true
autorestart=true
startsecs=15
stdout_logfile=/var/log/supervisor/lwn-simulator.log
stderr_logfile=/var/log/supervisor/lwn-simulator_err.log
directory=/LWN-Simulator
